# images/php7/config/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.5TTt7q/crontab installed on Wed Nov  1 16:00:05 2017)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)

*/10 * * * * /bin/date 1>> /data/cron.test.log 2>&1
