# images/sphinx/crontabs/root
* * * * * /root/sh/gen-delta-index >> /data/sphinx-delta.log 2>&1
